# ci/docker/rootfs/etc/services.d/api/run
#!/usr/bin/with-contenv sh
cd /app/
dotnet InstaBlaster.dll
